# final/roles/nginx/files/rsync_www.timer
[Unit]
Description=Rsync_WWW

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target