# final/roles/nginx/files/rsync_www.timer
[Unit]
Description=Rsync_WWW

[Timer]
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target