# final/roles/borg/files/borg-backup.timer
[Unit]
Description=Automated Borg Backup Timer

[Timer]
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target